# launchd/com.liltony.daytime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.liltony.daytime</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>/Users/tonybot/Lil-Tony-Trading-Brain/daytime_runner.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tonybot/Lil-Tony-Trading-Brain</string>

  <!-- 7:00 AM CT every weekday (Mon–Fri) -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/tonybot/Lil-Tony-Trading-Brain/logs/daytime.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tonybot/Lil-Tony-Trading-Brain/logs/daytime.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TZ</key>
    <string>America/Chicago</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>

</dict>
</plist>
